# systemd/system/static-arp.service
[Unit]
Description=Load static ARP entries from /etc/ethers
# Chạy sau khi mạng đã được cấu hình
After=network.target

[Service]
Type=oneshot
# Lệnh arp -f sẽ đọc và áp dụng tất cả các entry trong file /etc/ethers
ExecStart=/sbin/arp -f /etc/ethers
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target